# cron/finance.crontab
# m h  dom mon dow   command
*/5 * * * * /opt/Finance/cron/poll_stock_prices.py
