# metrics/finance.crontab
# m h  dom mon dow   command
*/5 * * * * /opt/Finance/metrics/poll_stock_prices.py
3 * * * * /opt/Finance/metrics/upload_to_s3.py
3 19 * * * /opt/Finance/analysis/download_from_s3.py 
